# supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
numprocs=1

[supervisord]
nodaemon=true
directory=/var/www/
